# custom-reboot.service
[Unit]
Description=Reboot weekly

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl reboot

[Install]
WantedBy=default.target
